# rssfeedy-pi.service
[Unit]
Description=RSSfeedy-Pi AI News Pushover Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alicankiraz
WorkingDirectory=/home/alicankiraz/RSSfeedy-Pi
ExecStart=/usr/bin/node /home/alicankiraz/RSSfeedy-Pi/dist/index.js
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

EnvironmentFile=/home/alicankiraz/RSSfeedy-Pi/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rssfeedy-pi

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/alicankiraz/RSSfeedy-Pi/data

[Install]
WantedBy=multi-user.target
